# systemd/devdash.service
[Unit]
Description=DevDash — AI-Powered Developer Desk Companion
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/raspi-demo
Environment=DISPLAY=:0
ExecStart=/home/pi/raspi-demo/.venv/bin/python -m devdash
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# GPIO cleanup on stop
ExecStopPost=/usr/bin/python3 -c "import RPi.GPIO as GPIO; GPIO.cleanup()"

[Install]
WantedBy=multi-user.target
