# systemd/devdash.service
[Unit]
Description=DevDash — Voice-First AI Developer Companion
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=samuel
WorkingDirectory=/home/samuel/devdash
Environment=DISPLAY=:1
ExecStart=/home/samuel/devdash/.venv/bin/python -m devdash
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Audio access for USB microphone
SupplementaryGroups=audio

[Install]
WantedBy=multi-user.target
